# spectrum-web.service
[Unit]
Description=Legion Spectrum Lighting Web UI
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/legion-spectrum-control/spectrum-web.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
